# homeNetMon.service
[Unit]
Description=HomeNetMon - Home Network Monitoring Service
Documentation=https://github.com/your-repo/homeNetMon
After=network.target network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
WorkingDirectory=/home/kyle/ClaudeCode/HomeNetMon
Environment=PYTHONPATH=/home/kyle/ClaudeCode/HomeNetMon
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NETWORK_RANGE=192.168.86.0/24
Environment=DEBUG=false
ExecStart=/usr/bin/python3 /home/kyle/ClaudeCode/HomeNetMon/app.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30
KillMode=mixed

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homeNetMon

# Environment variables
Environment=DATABASE_URL=sqlite:////home/kyle/ClaudeCode/HomeNetMon/homeNetMon.db
Environment=HOST=127.0.0.1
Environment=PORT=5000
Environment=DEBUG=false
Environment=NETWORK_RANGE=192.168.86.0/24
Environment=PING_INTERVAL=30
Environment=SCAN_INTERVAL=300
Environment=PING_TIMEOUT=3.0
Environment=MAX_WORKERS=50
Environment=DATA_RETENTION_DAYS=30

[Install]
WantedBy=multi-user.target