# systemd/homenetmon.service
[Unit]
Description=HomeNetMon Network Monitoring Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=homenetmon
Group=homenetmon
WorkingDirectory=/opt/homenetmon
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/homenetmon
Environment=DATABASE_URL=sqlite:////opt/homenetmon/data/homeNetMon.db
ExecStart=/opt/homenetmon/venv/bin/python app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homenetmon

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/homenetmon/data /opt/homenetmon/logs
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
